# midonet-cluster/src/rhel7/control/prerm
#!/bin/sh

if [ $1 -eq 0 ]; then
        # Package removal, not upgrade
        systemctl --no-reload disable midonet-cluster.service || true
        systemctl stop midonet-cluster.service || true
fi
