# midolman/src/rhel7/control/prerm
#!/bin/sh

if [ $1 -eq 0 ]; then
        # Package removal, not upgrade
        systemctl --no-reload disable midolman.service || true
        systemctl stop midolman.service || true
fi
