# midolmanj/src/deb/control/postrm
#!/bin/sh

# Midokura

if [ "$1" = purge ]; then
    update-rc.d midolmanj remove
    /usr/sbin/userdel midonet 2> /dev/null >/dev/null
fi

exit 0
