# midonet-cluster/src/deb/control/postinst
#!/bin/sh

NAME="midonet-cluster"
DATADIR=/var/lib/midonet-cluster
POST_DIR=/var/lib/dpkg/info

install -d /var/log/$NAME $DATADIR


test -e /lib/init/upstart-job \
    && ln -sf /lib/init/upstart-job /etc/init.d/midonet-cluster \
    || echo "upstart not installed"

systemctl preset midonet-cluster.service > /dev/null 2>&1 \
    || echo "systemd not installed"
systemctl enable midonet-cluster.service > /dev/null 2>&1 \
    || echo "systemd not installed"
