# midolman/src/deb/control/postinst
#!/bin/sh

NAME="midolman"
LOGDIR=/var/log/$NAME
DATADIR=/var/lib/$NAME
DBDIR=/var/db/$NAME
POST_DIR=/var/lib/dpkg/info

install -d $LOGDIR $DATADIR $DBDIR

mkdir -p /var/log/$NAME/minions

test -e /lib/init/upstart-job \
    && ln -sf /lib/init/upstart-job /etc/init.d/midolman \
    || echo "upstart not installed"

systemctl preset midolman.service > /dev/null 2>&1 \
    || echo "systemd not installed"
systemctl enable midolman.service > /dev/null 2>&1 \
    || echo "systemd not installed"
